# systemd/veris-backfill.service
[Unit]
Description=Veris pipeline — Stage 6: Backfill published_at (chain end)
After=network.target

[Service]
Type=oneshot
User=veris
WorkingDirectory=/home/veris/projects/veris
Environment=PATH=/home/veris/projects/veris/venv/bin
ExecStart=/home/veris/projects/veris/venv/bin/python3 -m stages.backfill
# No Restart= — oneshot services aren't restarted by systemd;
# if a stage fails, the next 3h fetch timer reruns the whole chain.

[Install]
WantedBy=multi-user.target
